# scripts/sportsbar-controller.service
[Unit]
Description=Sports Bar TV Controller
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/Sports-Bar-TV-Controller
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=%h/Sports-Bar-TV-Controller
Environment=FLASK_SECRET_KEY=sports-bar-production-key-2024
Environment=GIT_EXEC_PATH=/usr/bin
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
