# scripts/sportsbar-controller.service
[Unit]
Description=Sports Bar TV Controller
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/sportsbar-controller
Environment=PATH=/opt/sportsbar-controller/venv/bin
ExecStart=/opt/sportsbar-controller/venv/bin/python main.py --host 0.0.0.0 --port 5000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sportsbar-controller

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/sportsbar-controller

[Install]
WantedBy=multi-user.target
